# DEBIAN/postinst
#!/bin/bash
set -e

# =============================================
# IMMICH SMART INSTALLER - POSTINST संस्करण
# लेखक: बिजय पोख्रेल
# 📞 +977 9851125758  |  📧 merobiju@gmail.com
# =============================================

RED='\033[0;31m'
GREEN='\033[0;32m'
YELLOW='\033[1;33m'
BLUE='\033[0;34m'
NC='\033[0m'

clear
echo -e "${BLUE}════════════════════════════════════════════════════════${NC}"
echo -e "${GREEN}          IMMICH स्मार्ट इन्स्टलर${NC}"
echo -e "${GREEN}          लेखक: बिजय पोख्रेल${NC}"
echo -e "${YELLOW}          📞 +977 9851125758${NC}"
echo -e "${YELLOW}          📧 merobiju@gmail.com${NC}"
echo -e "${BLUE}════════════════════════════════════════════════════════${NC}"
echo ""

# =============================================
# DOCKER सेवा जाँच र सुरु
# =============================================
echo -e "${YELLOW}🔧 Docker सेवा जाँच गर्दै...${NC}"
if systemctl is-active --quiet docker; then
    echo -e "${GREEN}✓ Docker सेवा पहिले नै चलिरहेको छ${NC}"
else
    echo -e "${YELLOW}⚠ Docker सेवा सुरु गर्दै...${NC}"
    systemctl start docker
    systemctl enable docker
    echo -e "${GREEN}✓ Docker सेवा सुरु भयो${NC}"
fi

# Docker कमाण्ड चल्ने सुनिश्चित
if ! command -v docker &> /dev/null; then
    echo -e "${RED}❌ Docker फेला परेन। कृपया 'docker-ce' प्याकेज जाँच गर्नुहोस्।${NC}"
    exit 1
fi

if ! docker compose version &> /dev/null; then
    echo -e "${RED}❌ Docker Compose फेला परेन। कृपया 'docker-compose-plugin' प्याकेज जाँच गर्नुहोस्।${NC}"
    exit 1
fi

echo -e "${GREEN}✓ Docker र Compose तयार छन्${NC}"

# =============================================
# अपरेटिङ सिस्टम पहिचान
# =============================================
if ! command -v lsb_release &> /dev/null; then
    echo -e "${RED}❌ यो इन्स्टलर Ubuntu/Debian मा मात्र चल्छ${NC}"
    exit 1
fi

OS=$(lsb_release -is)
ARCH=$(uname -m)

if [[ "$OS" != "Ubuntu" && "$OS" != "Debian" ]]; then
    echo -e "${RED}❌ त्रुटि: यो इन्स्टलर Ubuntu/Debian मा मात्र काम गर्छ${NC}"
    echo -e "${YELLOW}   तपाईंको OS: $OS${NC}"
    exit 1
fi

if [[ "$ARCH" != "x86_64" && "$ARCH" != "aarch64" ]]; then
    echo -e "${RED}❌ त्रुटि: ६४-बिट आर्किटेक्चर आवश्यक छ${NC}"
    echo -e "${YELLOW}   तपाईंको आर्किटेक्चर: $ARCH${NC}"
    exit 1
fi

echo -e "${GREEN}✓ OS: $OS${NC}"
echo -e "${GREEN}✓ आर्किटेक्चर: $ARCH${NC}"

# =============================================
# ब्राउजर जाँच (स्वचालित)
# =============================================
if ! command -v firefox &> /dev/null && ! command -v google-chrome &> /dev/null && ! command -v chromium-browser &> /dev/null; then
    echo -e "${YELLOW}⚠ कुनै ब्राउजर फेला परेन। Firefox स्थापना गर्दै...${NC}"
    apt update && apt install -y firefox
fi

# =============================================
# RAM पत्ता लगाउने र ट्युनिङ
# =============================================
TOTAL_RAM_GB=$(free -g | awk '/^Mem:/{print $2}')

if [ "$TOTAL_RAM_GB" -lt 4 ]; then
    echo -e "${RED}❌ त्रुटि: न्यूनतम ४GB RAM आवश्यक छ${NC}"
    echo -e "${RED}   तपाईंको RAM: ${TOTAL_RAM_GB}GB${NC}"
    exit 1
elif [ "$TOTAL_RAM_GB" -lt 8 ]; then
    ENABLE_ML="false"
    SERVER_MEM="1G"
    WORKER_THREADS="1"
    echo -e "${YELLOW}⚠ कम RAM मोड (${TOTAL_RAM_GB}GB) - ML निष्क्रिय${NC}"
elif [ "$TOTAL_RAM_GB" -lt 16 ]; then
    ENABLE_ML="true"
    SERVER_MEM="2G"
    WORKER_THREADS="2"
    echo -e "${GREEN}✓ सामान्य मोड (${TOTAL_RAM_GB}GB) - ML सक्रिय${NC}"
else
    ENABLE_ML="true"
    SERVER_MEM="4G"
    WORKER_THREADS="4"
    echo -e "${GREEN}✓ उच्च प्रदर्शन मोड (${TOTAL_RAM_GB}GB) - पूर्ण सुविधाहरू${NC}"
fi

# =============================================
# पोर्ट पत्ता लगाउने (स्वचालित)
# =============================================
PORT=2283
while ss -tuln | grep -q ":$PORT "; do
    PORT=$((PORT + 1))
done
echo -e "${GREEN}✓ उपलब्ध पोर्ट: $PORT${NC}"

# =============================================
# IMMICH सेटअप
# =============================================
echo -e "${YELLOW}🔧 Immich सेटअप गर्दै...${NC}"

mkdir -p /opt/immich/{upload,pgdata,model-cache}
chmod 755 /opt/immich/upload

# सुरक्षित प्रमाणहरू उत्पन्न गर्ने
DB_PASSWORD=$(openssl rand -base64 32 | tr -d "=+/" | cut -c1-32)
JWT_SECRET=$(openssl rand -base64 32 | tr -d "=+/" | cut -c1-32)

# .env फाइल बनाउने
cat > /opt/immich/.env << ENVEOF
DB_PASSWORD=$DB_PASSWORD
DB_USERNAME=immich
DB_DATABASE_NAME=immich
JWT_SECRET=$JWT_SECRET
ENABLE_ML=$ENABLE_ML
WORKER_THREADS=$WORKER_THREADS
ENVEOF

# docker-compose.yml फाइल बनाउने
cat > /opt/immich/docker-compose.yml << 'DOCEOF'
version: '3.8'
name: immich
services:
  immich-server:
    image: ghcr.io/immich-app/immich-server:release
    container_name: immich_server
    ports:
      - "${PORT}:2283"
    volumes:
      - /opt/immich/upload:/usr/src/app/upload
      - /etc/localtime:/etc/localtime:ro
    env_file:
      - .env
    environment:
      - DB_PASSWORD=${DB_PASSWORD}
      - DB_USERNAME=${DB_USERNAME}
      - DB_DATABASE_NAME=${DB_DATABASE_NAME}
      - JWT_SECRET=${JWT_SECRET}
      - WORKER_THREADS=${WORKER_THREADS}
    depends_on:
      - database
      - redis
    restart: always

  immich-microservices:
    image: ghcr.io/immich-app/immich-microservices:release
    container_name: immich_microservices
    volumes:
      - /opt/immich/upload:/usr/src/app/upload
      - /etc/localtime:/etc/localtime:ro
    env_file:
      - .env
    environment:
      - DB_PASSWORD=${DB_PASSWORD}
      - DB_USERNAME=${DB_USERNAME}
      - DB_DATABASE_NAME=${DB_DATABASE_NAME}
      - JWT_SECRET=${JWT_SECRET}
    depends_on:
      - database
      - redis
    restart: always

  immich-machine-learning:
    image: ghcr.io/immich-app/immich-machine-learning:release
    container_name: immich_ml
    volumes:
      - /opt/immich/model-cache:/cache
    env_file:
      - .env
    restart: always
    profiles: ["ml"]

  database:
    image: tensorchord/pgvecto-rs:pg14-v0.2.0
    container_name: immich_postgres
    environment:
      POSTGRES_PASSWORD: ${DB_PASSWORD}
      POSTGRES_USER: immich
      POSTGRES_DB: immich
    volumes:
      - /opt/immich/pgdata:/var/lib/postgresql/data
    restart: always

  redis:
    image: redis:7.2-alpine
    container_name: immich_redis
    command: redis-server --appendonly yes --maxmemory 256mb --maxmemory-policy allkeys-lru
    restart: always
DOCEOF

# पोर्ट प्लेसहोल्डर प्रतिस्थापन गर्ने
sed -i "s/\${PORT}/$PORT/g" /opt/immich/docker-compose.yml

# =============================================
# IMMICH सुरु गर्ने
# =============================================
cd /opt/immich
echo -e "${YELLOW}🐳 Docker इमेजहरू तान्दै...${NC}"
docker compose pull

echo -e "${YELLOW}🚀 Immich सुरु गर्दै...${NC}"
if [ "$ENABLE_ML" = "true" ]; then
    docker compose --profile ml up -d
else
    docker compose up -d
fi

# =============================================
# TAILSCALE सेटअप स्क्रिप्ट
# =============================================
if ! command -v tailscale &> /dev/null; then
    echo -e "${YELLOW}🔒 Tailscale स्थापना गर्दै...${NC}"
    curl -fsSL https://tailscale.com/install.sh | sh
    echo -e "${GREEN}✓ Tailscale स्थापित भयो${NC}"
fi

# Tailscale सेटअप स्क्रिप्ट (पछि म्यानुअल रूपमा चलाउन)
cat > /usr/local/bin/immich-tailscale-setup << 'TSEOF'
#!/bin/bash
echo "════════════════════════════════════════════════════════"
echo "     IMMICH को लागि TAILSCALE सेटअप"
echo "     लेखक: बिजय पोख्रेल | 📞 +977 9851125758"
echo "════════════════════════════════════════════════════════"
echo ""
echo "🔐 Tailscale मा लगइन गर्दै (एक पटक मात्र)"
echo "   लिंक ब्राउजरमा खुल्नेछ"
echo ""
tailscale up
echo ""
echo "✅ Tailscale जडान भयो!"
echo "   तपाईं अब जहाँबाट पनि Immich पहुँच गर्न सक्नुहुन्छ"
echo ""
TSEOF
chmod +x /usr/local/bin/immich-tailscale-setup

# =============================================
# डेस्कटप सर्टकटहरू
# =============================================
echo -e "${YELLOW}📁 डेस्कटप सर्टकटहरू बनाउँदै...${NC}"

# Immich सर्टकट
cat > ~/Desktop/Immich-Photos.desktop << DESKEOF
[Desktop Entry]
Name=🖼️ Immich Photos
Comment=स्मार्ट फोटो ब्याकअप - बिजय पोख्रेल
Exec=xdg-open http://localhost:$PORT
Icon=immich
Terminal=false
Type=Application
Categories=Network;Photography;
StartupNotify=true
DESKEOF
chmod +x ~/Desktop/Immich-Photos.desktop

# Tailscale सर्टकट
cat > ~/Desktop/Tailscale-Setup.desktop << TAILDESK
[Desktop Entry]
Name=🔒 Tailscale सेटअप
Comment=टाढाबाट पहुँचको लागि Tailscale सेटअप (एक पटक)
Exec=gnome-terminal -- bash /usr/local/bin/immich-tailscale-setup
Icon=tailscale
Terminal=false
Type=Application
Categories=Network;
TAILDESK
chmod +x ~/Desktop/Tailscale-Setup.desktop

# प्रणाली अनुप्रयोगहरूमा प्रतिलिपि गर्ने
cp ~/Desktop/Immich-Photos.desktop /usr/share/applications/ 2>/dev/null || true

# =============================================
# आइकन बनाउने (Base64 बाट)
# =============================================
echo "iVBORw0KGgoAAAANSUhEUgAAAEAAAABACAYAAACqaXHeAAABhGlDQ1BJQ0MgcHJvZmlsZQAAKJF9kT1Iw0AcxV9TpUUqDnYQcchQnSyIijhKFYtgobQVWnUwufQLmjQkKS6OgmvBwY/FqoOLs64OroIg+AHi5uak6CIl/i8ptIjx4Lgf7+497t4BQqPCNLNrAtB020wl4mImuyoGX9GHMYwjJpmZJ8xKSo6n+3iHh3cX7lne53wTAUrCMsE3iGeY3UXDeJ7YxNVKq/icOMNqSpFvjGfK/LI84bIf+9h2kTMd5JqYJTKZFj2WHa+0eI5YMzVV+wWxj6iUaH5T5Tlnu1zxnM+8u9xqzb2t9k7uDqA/mi0vs5wGokQSSYQoqKCOCtJwUGWgKqHOuv0fN/xi+SRyVeVY4yBX1UGW7BAP+geP5qolI+O42A4IuH62/ZcZILwJXGsN23tt2PbnAJh9Aq45y19v2M1PkPQdXRyCvwPupS67aAeC+0DvtWHNXaTPQ6qUIs8DwMcZfVIGJG6DnjWn997O9h7+nWl1vwfWsnLPbXa8RAAAAAlwSFlzAAAuIwAALiMBeKU/dgAAAAd0SU1FB+cFFQURExQnNL4AAAAZdEVYdENvbW1lbnQAQ3JlYXRlZCB3aXRoIEdJTVBXgQ4XAAAABnRFWHRTb2Z0d2FyZQBnbm9tZS1zY3JlZW5zaG907wO/PgAABRZJREFUeJztW2tsG1UVPueO13bsxEnsOJXjJM7LcRraFGghPCoVKqhIfhQhUEWlokIVhPgBQhW/ECoBAiEhUUFoN4AiVKBQKXSp0gJKE9JSKKWlD6Bp0jRp4th5eG3H8dqzvfcefzizydpOYjvrbijfyUoe7575znfPmXPmzJk7QMBgMBgMBoPBYDAYDAaDwWAwGAwGg8FgMCoEf1UnYBSYf1fG01RlC9ABz1lPA7pS95fWp/AzNE6jSaAqXgXm9Y0btm/6AeB/sCj4MyQRxNOoKPRwXAwAQNndjz8/8v0VYNEkBMwUuK3rU8DdEy9kH/Ye2P3zN5gQMIvgnp0DAMDA+O8Uqll1iQpAMQhA1YzA0knnY5USyYYoyhTHnADzT3KOWQCMnxVYjA4gFypHvJtHvuKGRVh3EQCm/SvqAMRrs8B2kHlMCgBixWXAIuSiMTIYZtFiZgDYOZPvIgSJzH+cy4PNAFAWyojYbQBSSxZA7DsDQLwE6H1eAObUbyU/xAphBkDGzK8S/r/vBGwCwMxW2ASpCYCc9+byYTYBYGorDABq4h9iBfD5agVQYp3C5wLQ1QQAInZc/oIoH2YTAAbrdYBMOqyUcwnVQr6qAGgAWkGBMFmrAeRU2RmqZkLQDMBpAHDKYakDUMt4fA6oWQBMQ84KZkHVBcAtBoDqXK6qAOCWdQNQrd1UawCqOgBQy3nDcstW2ylcax4AoOx8AGj8i1yBdUgNwN4Eo1UO7gMAyHYJCO0CSJKjrmZ3AOxoFkBZt4yENy0AQ9QSAHkMgcTv4oj4BZRgHpPDAAArOwQQZAsoHQhC8Hm1ms2jD0RACQAkRwySqFOXZRBCkHkcsj+BsmY5Uxq8QyBpCIhTIEkEVH0AhA2HQRyHmikBwV88AHllF1RVUTvJ4cA8A4A9yZgSADjOAVGUAQDB5wXkMaCOu8Lw+u6bq9lixVrgqg4A+b9/kSgEKh0EQQiWAMhxYFoAckUwLQDJH4YkCSAIwDwGXj1pRPTALdVsseIAEEUSJFEEngfHj4OQoQwDgKDZAHw+L2RRAklwQPViINWyVV05ICpBECFJIoiiIHgOVBEA4QO3Qgz4TGlDpL/+eS1bo2oCEBWVIAAgC0FwvARBFIDISV/8W9eS25rNkSqyKqKpX7k/uvxE8iB5hhBc7QAgJH/5s+rX5/x4T7O5UlUjYc7k/TMfs+4FgGk8lCxIdHhC+YGSrm2/qH5t2p2zzXILqAgRAln4IIA2BwAMgAptL/v/Y5c/O+yp6oWgllmh3qT9/wJADgPAqMB0IDYvI13I9syuU67+W8uWYVUHAKkAIMhMADAqMAMIRuE4GCQCAIs7Z9nUcA8A4FgAGBWYAYSmGgHwWgBwLACMCswAIicEwG8CwLEAMCowA0jQAgC/DQDHAsCowAwgcUIA/FYAHAsAowIzgAQC8CgA+IUlRABgVGAGEKcFwK8FwLEAMCowA4hJAtBcBqC4DQBWB4BRgRlAjBEATQkAZgFgVGAEEDgOAKAAIJcBwLAgA4AQEwCUBYCMsSALgBATAJQFgAABtH4CQLhQ4EQBWAByIWQFgLIBYDgWAHYFMqEoAOT/wecCYLgUAELWCqjFgJSD4Vp1IQDwCjMA9ZCPCQDyWVg2AGwLZFQWAPlPWLKqSgCYcDBMAyCvYVlVFwA2GcxYAKgQ7Kg0AGw2mKkAMGxgWQAw3Hh4AhMDYFjNtAAwjOFTIDMAsF1gaoFhNtMCwDCGqQUmBcAwgokFJgXAMIqJBWYJQzMJTDsYxi4wLQAMzQYjAdAaPHgCc52AwyFgMvha2i3x4TPwP8H4Bw1Bc7aEYh12AAAAAElFTkSuQmCC" | base64 -d | tee /usr/share/pixmaps/immich.png > /dev/null

# =============================================
# प्रमाणहरू सुरक्षित गर्ने
# =============================================
cat > /root/.immich_credentials << CREDEOF
════════════════════════════════════════════════════════
     IMMICH स्थापना - बिजय पोख्रेल
     📞 +977 9851125758 | 📧 merobiju@gmail.com
════════════════════════════════════════════════════════
स्थापना मिति: $(date)
प्रणाली: ${TOTAL_RAM_GB}GB RAM
पहुँच URL: http://localhost:$PORT

डाटाबेस प्रमाणहरू:
  प्रयोगकर्ता नाम: immich
  पासवर्ड: $DB_PASSWORD
  डाटाबेस: immich

JWT गोप्य: $JWT_SECRET

आदेशहरू:
  लग हेर्ने: cd /opt/immich && docker compose logs -f
  पुन: सुरु गर्ने: cd /opt/immich && docker compose restart
  रोक्ने: cd /opt/immich && docker compose down
  अद्यावधिक: cd /opt/immich && docker compose pull && docker compose up -d
════════════════════════════════════════════════════════
CREDEOF
chmod 600 /root/.immich_credentials

# =============================================
# सेवाहरूको लागि पर्खने
# =============================================
echo -e "${YELLOW}⏳ Immich सुरु हुनको लागि पर्खँदै...${NC}"
sleep 10

# =============================================
# अन्तिम सफलता सन्देश
# =============================================
clear
echo -e "${BLUE}════════════════════════════════════════════════════════${NC}"
echo -e "${GREEN}                   ✅ स्थापना पूरा भयो! ✅${NC}"
echo -e "${BLUE}════════════════════════════════════════════════════════${NC}"
echo ""
echo -e "${GREEN}          IMMICH स्मार्ट इन्स्टलर${NC}"
echo -e "${GREEN}          लेखक: बिजय पोख्रेल${NC}"
echo -e "${YELLOW}          📞 +977 9851125758${NC}"
echo -e "${YELLOW}          📧 merobiju@gmail.com${NC}"
echo ""
echo -e "${BLUE}════════════════════════════════════════════════════════${NC}"
echo ""
echo -e "${GREEN}🔗 Immich पहुँच गर्नुहोस्:${NC} ${BLUE}http://localhost:$PORT${NC}"
echo -e "${GREEN}📁 अपलोडहरू:${NC} ${BLUE}/opt/immich/upload${NC}"
echo -e "${GREEN}🔐 प्रमाणहरू:${NC} ${BLUE}/root/.immich_credentials${NC}"
echo ""
echo -e "${YELLOW}🖼️  डेस्कटप सर्टकट: 'Immich Photos'${NC}"
echo -e "${YELLOW}🔒 टाढा पहुँच: डेस्कटप 'Tailscale सेटअप' → एक पटक चलाउनुहोस्${NC}"
echo ""
echo -e "${BLUE}════════════════════════════════════════════════════════${NC}"
echo -e "${GREEN}           🙏 स्मार्ट इन्स्टलर प्रयोग गर्नुभएकोमा धन्यवाद 🙏${NC}"
echo -e "${BLUE}════════════════════════════════════════════════════════${NC}"
echo -e "${YELLOW}      सहायताको लागि: +977 9851125758 | merobiju@gmail.com${NC}"
echo -e "${BLUE}════════════════════════════════════════════════════════${NC}"

# ब्राउजर स्वत: खोल्ने
sleep 2
xdg-open http://localhost:$PORT 2>/dev/null || true

exit 0
